# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/workspace/development/humpback-acoustic-embed/logs/supervisord.log
logfile_maxbytes=20MB
logfile_backups=10
pidfile=/tmp/supervisord.pid
childlogdir=/workspace/development/humpback-acoustic-embed/logs
nodaemon=false
minfds=1024

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[group:humpback]
programs=humpback-api,humpback-worker

; Build the frontend during deploy. The API serves the built SPA from
; src/humpback/static/dist/ in production.

[program:humpback-api]
directory=/workspace/development/humpback-acoustic-embed
command=/bin/bash -lc 'exec uv run humpback-api'
autostart=true
autorestart=true
startsecs=5
startretries=3
stopsignal=TERM
stopwaitsecs=30
priority=10
stdout_logfile=/workspace/development/humpback-acoustic-embed/logs/api.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile=/workspace/development/humpback-acoustic-embed/logs/api_error.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=10
stopasgroup=true
killasgroup=true

[program:humpback-worker]
directory=/workspace/development/humpback-acoustic-embed
command=/bin/bash -lc 'exec uv run humpback-worker'
autostart=true
autorestart=true
startsecs=5
startretries=3
stopsignal=TERM
stopwaitsecs=60
priority=20
stdout_logfile=/workspace/development/humpback-acoustic-embed/logs/worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile=/workspace/development/humpback-acoustic-embed/logs/worker_error.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=10
stopasgroup=true
killasgroup=true
